# deploy/callisto-jupiter.service
[Unit]
Description=Callisto Jupiter metrics exporter
Documentation=https://callistosignal.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Config is read from /etc/callisto-jupiter/config.toml (see config.example.toml).
# Optionally override via an environment file instead:
#   EnvironmentFile=-/etc/callisto-jupiter/env
ExecStart=/opt/callisto-jupiter/bin/callisto-jupiter
Restart=always
RestartSec=10
# Run unprivileged; disk_path must be readable by this user.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
# Persisted store-and-forward buffer lives in /var/lib/callisto-jupiter.
# StateDirectory makes it writable under DynamicUser + ProtectSystem=strict and
# exports $STATE_DIRECTORY, which the agent's default buffer path honors.
StateDirectory=callisto-jupiter

[Install]
WantedBy=multi-user.target
